# SmartHome/fritzbox2telegraf/crontab
*/5 * * * * /home/fritzbox/telegrafFritzBox.sh
